# deploy/pingwatch-kiosk-restart.service
[Unit]
Description=Restart pingwatch-kiosk service

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart pingwatch-kiosk.service
